# tigase-server.upstart

description "Starts the Tigase XMPP server"

start on pld.network-started
stop on pld.shutdown-started

respawn

console output

script
	export TIGASE_CONSOLE_LOG=/var/log/tigase-server/tigase-console.log
	exec /usr/bin/tigase-server run
end script

# vi: ft=upstart
